# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:tellstickmqtt]
command=java -Djava.security.egd=file:/dev/./urandom -jar /app.jar

[program:telldusd]
command=telldusd --nodaemon